# airootfs/etc/systemd/system/sockets.target.wants/cups.socket
[Unit]
Description=CUPS Scheduler
PartOf=cups.service

[Socket]
ListenStream=/run/cups/cups.sock
RemoveOnStop=on

[Install]
WantedBy=sockets.target
